# files/upstart/cloud-gateway.conf
description "Cloud gateway"

start on runlevel [23]
stop on shutdown

exec /usr/java/jdk1.7.0_71/bin/java -Dgateway.id=appgw01 -Xms2G -Xmx2G -jar /datayes/gateway/web/gateway.jar 8081 8444

respawn
